# deploy/systemd/savo_edge.service
[Unit]
Description=Savo Edge Service
After=network.target

[Service]
Type=simple
ExecStart=/usr/bin/env bash /home/ahnaf-tahmid/Savo_Pi/tools/scripts/run_edge.sh
Restart=on-failure

[Install]
WantedBy=multi-user.target
